# .config/systemd/user/tpm-fido.service
[Unit]
Description=tpm-fido - tpm2 as fido token

[Service]
Type=simple
ExecStart=-/bin/bash -c 'PATH=~/go/bin:$PATH tpm-fido'
Restart=no

[Install]
WantedBy=default.target
