# .config/containers/systemd/plantuml.container
[Unit]
Description=Local plantuml server

[Container]
Image=plantuml/plantuml-server:jetty
PublishPort=9742:8080
